# deploy/africa-frontend.service
[Unit]
Description=Africa Frontend
After=network.target

[Service]
User=deploy
WorkingDirectory=/opt/africa/africa-frontend
ExecStart=/usr/bin/node .next/standalone/server.js
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
